# supervisord.conf
[program:node]
command = nodemon
directory = /home/dotcloud/current
